# components/sentry/.docker/files/supervisor.d/redis.conf
[program:redis]
command=redis-server
autostart=true
autorestart=true
user=root
redirect_stderr=true
startsecs=0
