# components/sentry/.docker/files/supervisor.d/sentry-cron.conf
[program:sentry-cron]
command=sentry run cron
autostart=true
autorestart=true
user=root
redirect_stderr=true
startsecs=0
