# components/sentry/.docker/files/supervisor.d/sentry-worker.conf
[program:sentry-worker]
command=sentry run worker
autostart=true
autorestart=true
user=root
redirect_stderr=true
startsecs=0
